# garmin/launchd/com.takt.garmin-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.takt.garmin-sync</string>

  <!-- Startet den Wrapper (setzt PATH, cd ins Repo, sync.py --days 3, Logging).
       Repo liegt außerhalb ~/Documents -> kein macOS-TCC, keine Sonderfreigabe nötig. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ericschmitt/Takt-Web/garmin/launchd/run-sync.sh</string>
  </array>

  <!-- 8x/Tag alle 3h ab 08:30: 08:30 / 11:30 / 14:30 / 17:30 / 20:30 / 23:30 / 02:30 / 05:30 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <!-- launchd-eigene Ausgaben (Wrapper-Start/-Fehler). Der eigentliche Sync-Log ist garmin/logs/sync.log. -->
  <key>StandardOutPath</key>
  <string>/Users/ericschmitt/Takt-Web/garmin/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ericschmitt/Takt-Web/garmin/logs/launchd.err.log</string>

  <!-- Beim Laden NICHT sofort feuern; nur zu den geplanten Zeiten. -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
